# shelley.service
[Unit]
Description=Shelley AI Assistant Service
Documentation=https://exe.dev
After=network.target
Requires=shelley.socket

[Service]
Type=simple
ExecStart=/usr/local/bin/shelley
Restart=always
RestartSec=5
User=exedev
Group=exedev
WorkingDirectory=/home/exedev

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
